# kahla-scanner/systemd/kahla-scanner.service
[Unit]
Description=Kahla Scanner
After=network.target

[Service]
Type=simple
User=scanner
WorkingDirectory=/opt/kahla-scanner
Environment="PATH=/opt/kahla-scanner/venv/bin"
EnvironmentFile=/opt/kahla-scanner/.env
ExecStart=/opt/kahla-scanner/venv/bin/python main.py
Restart=always
RestartSec=10
# Log to journalctl
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
